# deploy/nara001-backup.timer
[Unit]
Description=Daily Nara001 backup

[Timer]
OnCalendar=*-*-* 03:20:00
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
