# deploy/nara001-backup.service
[Unit]
Description=Back up Nara001 database and media

[Service]
Type=oneshot
EnvironmentFile=-/etc/nara001/backup.env
ExecStart=/usr/local/sbin/nara001-backup
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
