# deploy/nara001-cloudflare-allowlist.service
[Unit]
Description=Refresh nara001 Cloudflare origin allowlist
After=network-online.target nginx.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/nara001-refresh-cloudflare-allowlist
